# installers/linux/orchestrator.service
[Unit]
Description=Claude Code Orchestrator
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=orchestrator
WorkingDirectory=/home/orchestrator/claude-orchestrator
Environment=NODE_ENV=production
Environment=PATH=/usr/local/bin:/usr/bin:/bin
ExecStart=/usr/bin/node packages/backend/dist/server.js
Restart=always
RestartSec=5
TimeoutStopSec=20

[Install]
WantedBy=multi-user.target
